# bot3.conf
[program:ASCIIBot]
directory=/home/ubuntu/ASCIIBot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/ASCIIBot/logfile_err.log
stdout_logfile=/home/ubuntu/ASCIIBot/logfile.log